# scripts/com.bamboohr.ps-gift-card-incentives.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.bamboohr.ps-gift-card-incentives</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>REPO_PATH/scripts/refresh_and_push.sh</string>
  </array>

  <!-- Run at top of every hour -->
  <key>StartCalendarInterval</key>
  <dict>
    <key>Minute</key>
    <integer>0</integer>
  </dict>

  <key>StandardOutPath</key>
  <string>/Users/jyorgason/Library/Logs/ps-gift-card-incentives.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/jyorgason/Library/Logs/ps-gift-card-incentives.log</string>

  <key>RunAtLoad</key>
  <true/>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/usr/local/bin:/usr/bin:/bin:/opt/homebrew/bin</string>
    <key>HOME</key>
    <string>/Users/jyorgason</string>
  </dict>
</dict>
</plist>
